# examples/systemd/clar-focus-recover.service
[Unit]
Description=Recover Clar Focus session state on login
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/clar-focus focus recover

[Install]
WantedBy=default.target
